# debian/prerm
#!/bin/bash
# This script is run before the package is removed

set -e

# stop cgminer server
invoke-rc.d cgminer stop
